# host/handslap/systemd/my-secondly.service
[Unit]
Description=My every second task

[Service]
ExecStart=/root/sbin/my-secondly
