# host/handslap/systemd/my-secondly.timer
[Unit]
Description=Run my-secondly.service every second

[Timer]
OnBootSec=0
OnUnitActiveSec=5
AccuracySec=1

[Install]
WantedBy=timers.target
